# docker/modules/os/fs/etc/profile.d/alias.sh
#!/bin/bash

alias d='ls -F'
alias l.='ls -d .* --color=tty'
alias ll='ls -l -h --color=tty'
alias ls='ls --color=tty'
alias la='ls -laSh --color=tty'
alias s='cd ..'
alias vi='vim'